# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Scan with radios to determine noise to fix signal icons in luci
timeout 1 iwinfo wlan0 scan &
timeout 1 iwinfo wlan1 scan &

# Turn on button LED
echo 255 > /sys/class/leds/c2600:white:general/brightness

exit 0
